# resources/.config/profile.d/90-ippub.sh

# Print my public IP address
if command -v curl &>/dev/null
then
    ippub() {
        curl ifconfig.me 2> /dev/null
    }
fi
