# resources/.config/profile.d/90-rsync.sh

if command -v rsync &>/dev/null; then
    alias rsync="rsync --info=progress2"
fi
