# system/home/.config/systemd/user/easyeffects.service
[Unit]
Description=Start EasyEffects

[Service]
ExecStart=/usr/bin/easyeffects --gapplication-service
Restart=on-failure

[Install]
WantedBy=default.target
