# user/.config/systemd/user/easyeffects.service
[Unit]
Description=Start EasyEffects

[Service]
ExecStart=/usr/bin/easyeffects --gapplication-service
Restart=always

[Install]
WantedBy=default.target
